# dump-extensions/openpower-dumps/create_openpower_dump_dirs.service.in
[Unit]
Description=Create openpower dump folders
Before=xyz.openbmc_project.Dump.Manager.service

[Service]
ExecStart=/bin/sh -c 'mkdir -p {HOSTBOOT_DUMP_PATH}'
ExecStart=/bin/sh -c 'mkdir -p {HARDWARE_DUMP_PATH}'
Type=oneshot

[Install]
WantedBy=multi-user.target
